# backup/rest-server.service
[Unit]
Description=Restic REST Server (Append-Only)
After=network.target

[Service]
User=backup
Group=backup
ExecStart=/usr/local/bin/rest-server --path /srv/backup --append-only --no-auth --listen :8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
